# reflex-strategy-vps 2.service
[Unit]
Description=Reflex-Strategy VPS Brain
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/reflex-strategy-system
Environment="PYTHONPATH=/root/reflex-strategy-system"
Environment="VPS_HOST=0.0.0.0"
Environment="VPS_PORT=8001"
Environment="DB_HOST=localhost"
Environment="DB_PORT=5432"
Environment="DB_NAME=reflex_strategy"
Environment="DB_USER=postgres"
Environment="DB_PASSWORD=postgres"
Environment="LLM_API_KEY="
Environment="LLM_MODEL=gpt-4o-mini"
ExecStart=/usr/bin/python3 -m uvicorn vps.server:app --host 0.0.0.0 --port 8001
Restart=always
RestartSec=10
StandardOutput=append:/var/log/reflex-strategy-vps.log
StandardError=append:/var/log/reflex-strategy-vps-error.log

[Install]
WantedBy=multi-user.target
